# packaging/systemd/aargal.service
[Unit]
Description=Aargal - Intelligent Request Barrier
After=network.target

[Service]
Type=simple

User=aargal
Group=aargal

ExecStart=/usr/local/bin/aargal --config /etc/aargal/aargal.toml

Restart=on-failure
RestartSec=5

# Security hardening (safe defaults)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/nginx

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
